# deploy/systemd/planning-engine-api.service
[Unit]
Description=Planning Engine FastAPI Service
After=network.target

[Service]
Type=simple
User=planning-engine
Group=planning-engine
WorkingDirectory=/opt/planning-engine/repo
Environment="PATH=/opt/planning-engine/venv/bin"
Environment="PYTHONPATH=/opt/planning-engine/repo/src"
ExecStart=/opt/planning-engine/venv/bin/uvicorn apps.api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=planning-engine-api

[Install]
WantedBy=multi-user.target
